# solarflow-control/rootfs/etc/s6-overlay/s6-rc.d/solarflow-control/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start Solarflow Control
# ==============================================================================

bashio::log.info "Running Solarflow Control"

# if the fork was not selected, then run the original solarflow-control
if bashio::config.false 'use_fork'; then
    cd /solarflow
    exec python3 /solarflow/solarflow-control.py
else
    cd /solarflow_fork
    exec python3 /solarflow_fork/solarflow-control.py
fi
