# solarflow-control/rootfs/etc/s6-overlay/s6-rc.d/solarflow-control/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start Solarflow Control
# ==============================================================================

# if the fork was not selected, then run the original solarflow-control
if bashio::config.true 'use_fork'; then
    bashio::log.info "Running Solarflow Control Fork"
    cd /solarflow_fork || exit
    exec python3 /solarflow_fork/solarflow-control.py
else
    bashio::log.info "Running Solarflow Control"
    cd /solarflow || exit
    exec python3 /solarflow/solarflow-control.py
fi
