# solarflow-control/rootfs/etc/s6-overlay/s6-rc.d/solarflow-control/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start Solarflow Control
# ==============================================================================

bashio::log.info "Running Solarflow Control"
cd /solarflow
exec python3 /solarflow/solarflow-control.py
