# solarflow-control/rootfs/etc/s6-overlay/s6-rc.d/init-solarflow-control/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Solarflow Control
# Configures the config.ini for Solarflow Control
# ==============================================================================

declare dtu_type
declare use_cloud
declare grid_charge
declare smartmeter_type

declare product_id
declare device_id
declare full_charge_interval
declare control_bypass
declare control_soc
declare battery_low
declare battery_high
declare cloud_web_user
declare cloud_web_password
declare web_token_url
declare web_device_list_url
declare cloud_mqtt_port
declare cloud_mqtt_host
declare cloud_mqtt_user
declare cloud_mqtt_pwd
declare roofsolar_base_topic
declare mqtt_host
declare mqtt_port
declare mqtt_user
declare mqtt_password
declare dtu_base_topic
declare sf_inverter_channels
declare inverter_max_power

declare opendtu_inverter_serial
declare ahoydtu_inverter_id

declare sm_rapid_change_diff
declare sm_zero_offset
declare sm_scaling_factor
declare sm_generic_base_topic
declare sm_generic_cur_accessor
declare sm_generic_total_accessor

declare sm_poweropti_user
declare sm_poweropti_password

declare sm_shelly3em_base_topic

declare min_charge_power
declare min_discharge_power
declare max_discharge_power
declare limit_inverter
declare max_inverter_limit
declare min_inverter_limit
declare discharge_during_daytime
declare latitude
declare longitude
declare sunrise_offset
declare sunset_offset
declare battery_low
declare battery_high

declare min_grid_charge_power
declare max_grid_charge_power

bashio::log.info 'Updating Solarflow Control configuration...'

# Read values from the add-on configuration and update config.ini
dtu_type=$(bashio::config 'dtu_type')
sed -i "s|<dtu_type>|$dtu_type|g" /solarflow/config.ini

use_cloud=$(bashio::config 'use_cloud')
sed -i "s|<use_cloud>|$use_cloud|g" /solarflow/config.ini

grid_charge=$(bashio::config 'grid_charge')
sed -i "s|<grid_charge>|$grid_charge|g" /solarflow/config.ini

smartmeter_type=$(bashio::config 'smartmeter_type')
sed -i "s|<smartmeter_type>|$smartmeter_type|g" /solarflow/config.ini

product_id=$(bashio::config 'product_id')
sed -i "s|<product_id>|$product_id|g" /solarflow/config.ini

device_id=$(bashio::config 'device_id')
sed -i "s|<device_id>|$device_id|g" /solarflow/config.ini

full_charge_interval=$(bashio::config 'full_charge_interval')
sed -i "s|<full_charge_interval>|$full_charge_interval|g" /solarflow/config.ini

control_bypass=$(bashio::config 'control_bypass')
sed -i "s|<control_bypass>|$control_bypass|g" /solarflow/config.ini

control_soc=$(bashio::config 'control_soc')
sed -i "s|<control_soc>|$control_soc|g" /solarflow/config.ini

battery_low=$(bashio::config 'battery_low')
sed -i "s|<battery_low>|$battery_low|g" /solarflow/config.ini

battery_high=$(bashio::config 'battery_high')
sed -i "s|<battery_high>|$battery_high|g" /solarflow/config.ini

cloud_web_user=$(bashio::config 'cloud_web_user')
sed -i "s|<cloud_web_user>|$cloud_web_user|g" /solarflow/config.ini

cloud_web_password=$(bashio::config 'cloud_web_password')
sed -i "s|<cloud_web_password>|$cloud_web_password|g" /solarflow/config.ini

web_token_url=$(bashio::config 'web_token_url')
sed -i "s|<web_token_url>|$web_token_url|g" /solarflow/config.ini

web_device_list_url=$(bashio::config 'web_device_list_url')
sed -i "s|<web_device_list_url>|$web_device_list_url|g" /solarflow/config.ini

cloud_mqtt_port=$(bashio::config 'cloud_mqtt_port')
sed -i "s|<cloud_mqtt_port>|$cloud_mqtt_port|g" /solarflow/config.ini

cloud_mqtt_host=$(bashio::config 'cloud_mqtt_host')
sed -i "s|<cloud_mqtt_host>|$cloud_mqtt_host|g" /solarflow/config.ini

cloud_mqtt_user=$(bashio::config 'cloud_mqtt_user')
sed -i "s|<cloud_mqtt_user>|$cloud_mqtt_user|g" /solarflow/config.ini

cloud_mqtt_password=$(bashio::config 'cloud_mqtt_password')
sed -i "s|<cloud_mqtt_password>|$cloud_mqtt_password|g" /solarflow/config.ini

roofsolar_base_topic=$(bashio::config 'roofsolar_base_topic')
sed -i "s|<roofsolar_base_topic>|$roofsolar_base_topic|g" /solarflow/config.ini

mqtt_host=$(bashio::config 'mqtt_host')
sed -i "s|<mqtt_host>|$mqtt_host|g" /solarflow/config.ini

mqtt_port=$(bashio::config 'mqtt_port' '1883')
sed -i "s|<mqtt_port>|$mqtt_port|g" /solarflow/config.ini

mqtt_user=$(bashio::config 'mqtt_user')
sed -i "s|<mqtt_user>|$mqtt_user|g" /solarflow/config.ini

mqtt_password=$(bashio::config 'mqtt_password')
sed -i "s|<mqtt_password>|$mqtt_password|g" /solarflow/config.ini

dtu_base_topic=$(bashio::config 'dtu_base_topic')
sed -i "s|<dtu_base_topic>|$dtu_base_topic|g" /solarflow/config.ini

sf_inverter_channels=$(bashio::config 'sf_inverter_channels')
sed -i "s|<sf_inverter_channels>|$sf_inverter_channels|g" /solarflow/config.ini

inverter_max_power=$(bashio::config 'inverter_max_power')
sed -i "s|<inverter_max_power>|$inverter_max_power|g" /solarflow/config.ini

opendtu_inverter_serial=$(bashio::config 'opendtu_inverter_serial')
sed -i "s|<opendtu_inverter_serial>|$opendtu_inverter_serial|g" /solarflow/config.ini

ahoydtu_inverter_id=$(bashio::config 'ahoydtu_inverter_id')
sed -i "s|<ahoydtu_inverter_id>|$ahoydtu_inverter_id|g" /solarflow/config.ini

sm_rapid_change_diff=$(bashio::config 'sm_rapid_change_diff')
sed -i "s|<sm_rapid_change_diff>|$sm_rapid_change_diff|g" /solarflow/config.ini

sm_zero_offset=$(bashio::config 'sm_zero_offset')
sed -i "s|<sm_zero_offset>|$sm_zero_offset|g" /solarflow/config.ini

sm_scaling_factor=$(bashio::config 'sm_scaling_factor')
sed -i "s|<sm_scaling_factor>|$sm_scaling_factor|g" /solarflow/config.ini

sm_generic_base_topic=$(bashio::config 'sm_generic_base_topic')
sed -i "s|<sm_generic_base_topic>|$sm_generic_base_topic|g" /solarflow/config.ini

sm_generic_cur_accessor=$(bashio::config 'sm_generic_cur_accessor' ' ')
sed -i "s|<sm_generic_cur_accessor>|$sm_generic_cur_accessor|g" /solarflow/config.ini

sm_generic_total_accessor=$(bashio::config 'sm_generic_total_accessor' ' ')
sed -i "s|<sm_generic_total_accessor>|$sm_generic_total_accessor|g" /solarflow/config.ini

sm_poweropti_user=$(bashio::config 'sm_poweropti_user' ' ')
sed -i "s|<sm_poweropti_user>|$sm_poweropti_user|g" /solarflow/config.ini

sm_poweropti_password=$(bashio::config 'sm_poweropti_password' ' ')
sed -i "s|<sm_poweropti_password>|$sm_poweropti_password|g" /solarflow/config.ini

sm_shelly3em_base_topic=$(bashio::config 'sm_shelly3em_base_topic')
sed -i "s|<sm_shelly3em_base_topic>|$sm_shelly3em_base_topic|g" /solarflow/config.ini

min_charge_power=$(bashio::config 'min_charge_power')
sed -i "s|<min_charge_power>|$min_charge_power|g" /solarflow/config.ini

max_discharge_power=$(bashio::config 'max_discharge_power')
sed -i "s|<max_discharge_power>|$max_discharge_power|g" /solarflow/config.ini

limit_inverter=$(bashio::config 'limit_inverter')
sed -i "s|<limit_inverter>|$limit_inverter|g" /solarflow/config.ini

max_inverter_limit=$(bashio::config 'max_inverter_limit')
sed -i "s|<max_inverter_limit>|$max_inverter_limit|g" /solarflow/config.ini

min_inverter_limit=$(bashio::config 'min_inverter_limit')
sed -i "s|<min_inverter_limit>|$min_inverter_limit|g" /solarflow/config.ini

discharge_during_daytime=$(bashio::config 'discharge_during_daytime')
sed -i "s|<discharge_during_daytime>|$discharge_during_daytime|g" /solarflow/config.ini

latitude=$(bashio::config 'latitude' ' ')
sed -i "s|<latitude>|$latitude|g" /solarflow/config.ini

longitude=$(bashio::config 'longitude' ' ')
sed -i "s|<longitude>|$longitude|g" /solarflow/config.ini

sunrise_offset=$(bashio::config 'sunrise_offset')
sed -i "s|<sunrise_offset>|$sunrise_offset|g" /solarflow/config.ini

sunset_offset=$(bashio::config 'sunset_offset')
sed -i "s|<sunset_offset>|$sunset_offset|g" /solarflow/config.ini

battery_low=$(bashio::config 'battery_low')
sed -i "s|<battery_low>|$battery_low|g" /solarflow/config.ini

battery_high=$(bashio::config 'battery_high')
sed -i "s|<battery_high>|$battery_high|g" /solarflow/config.ini

min_grid_charge_power=$(bashio::config 'min_grid_charge_power')
sed -i "s|<min_grid_charge_power>|$min_grid_charge_power|g" /solarflow/config.ini

max_grid_charge_power=$(bashio::config 'max_grid_charge_power')
sed -i "s|<max_grid_charge_power>|$max_grid_charge_power|g" /solarflow/config.ini

min_discharge_power=$(bashio::config 'min_discharge_power')
sed -i "s|<min_discharge_power>|$min_discharge_power|g" /solarflow/config.ini